# contrib/initscripts/systemd/xpd.service
[Unit]
Description=XP daemon serivce
After=network.target

[Service]
Type=simple
User=XP
ExecStart=/opt/XP/XPd

[Install]
WantedBy=multi-user.target
